# flask/scripts/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=./supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=./supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=quasiben                                     ; default user
childlogdir=/home/bzaitlen/image_solr/scripts/supervisord_child_logs               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:image_solr]
command=/home/bzaitlen/miniconda/envs/image_solr/bin/python -u  /home/bzaitlen/miniconda/envs/image_solr/bin/gunicorn -c /home/bzaitlen/image_solr/scripts/gunicorn_cfg.py --pythonpath=/home/bzaitlen/image_solr server:app
stdout_logfile = /home/bzaitlen/image_solr/scripts/supervisord_child_logs/image_solr.log
redirect_stderr = true
autostart = true
killasgroup = true
